# templates/etc/profile.d/simp.csh.erb
# This file managed by Puppet.
# Any changes will be overwritten at the next run.

<% if not @user_whitelist.empty? then -%>
foreach user (<%= @user_whitelist.join(' ') %>)
  if ( "$user" == "$USER" ) then
    exit
  endif
end

<% end -%>
<% if @prepend['csh'] then -%>
<%= @prepend['csh'] %>

<% end -%>
set autologout=<%= @session_timeout %>
umask <%= @umask %>
tty -s
if ( $? == 0 ) mesg <%= @mesg ? 'y' : 'n' %>
<% if @append['csh'] then -%>

<%= @append['csh'] %>
<% end -%>
